# test/rootfs-tweaks/etc/init/superman.conf
# superman - configure and load superman
#
# This task loads the superman kernel module and daemon

description	"configure and load SUPERMAN"

start on  local-filesystems

exec /bin/configure_node
