# test/rootfs-tweaks/lib/systemd/system/superman.service
[Unit]
Description=configure and load SUPERMAN

[Service]
Type=simple
ExecStart=/bin/configure_node

[Install]
WantedBy=multi-user.target
